# deploy/temms.service
[Unit]
Description=TEMMS - Tactical Edge Model Management System
Documentation=https://github.com/yourusername/temms
After=network.target

[Service]
Type=simple
User=temms
Group=temms
WorkingDirectory=/opt/temms
Environment="PATH=/opt/temms/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/etc/temms/temms.env
ExecStart=/opt/temms/venv/bin/temms daemon start --foreground --config /etc/temms/temms.yaml --host ${TEMMS_HOST} --port ${TEMMS_PORT}
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Basic hardening. TEMMS still needs write access to /var/lib/temms and /var/log/temms.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/temms /var/log/temms /etc/temms

# Environment:
# - TEMMS_API_TOKEN protects control, Hub Lite, and Web UI write endpoints.
# - TEMMS_HUB_URL enables optional online Hub Lite enrollment and heartbeat sync.

[Install]
WantedBy=multi-user.target
